# collector/systemd/renogymon-bms-collector.service
[Unit]
Description=Renogy BMS Collector
After=bluetooth.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=renogymon-collector
Group=renogymon-collector
SupplementaryGroups=dialout
EnvironmentFile=-/etc/default/renogymon-collector
ExecStart=/usr/bin/renogymon-bms-collector ${COLLECTOR_ARGS}
Restart=on-failure
RestartSec=30
Environment=RUST_LOG=info
Environment=COLLECTOR_ARGS=bt2
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
