# archiver/systemd/renogymon-archiver-export.timer
[Unit]
Description=Run Renogy metrics export daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
