# archiver/systemd/renogymon-archiver-export.service
[Unit]
Description=Export Renogy metrics to Parquet
After=victoria-metrics.service

[Service]
Type=oneshot
User=renogymon-archiver
Group=renogymon-archiver
StateDirectory=renogymon-archiver
ExecStart=/usr/bin/renogymon-archiver export
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
