# puller/systemd/renogy-archiver-puller.timer
[Unit]
Description=Pull Renogy Parquet archives daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
