# aprs/systemd/renogymon-aprs.service
[Unit]
Description=Renogy APRS Reporter
# gpsd ordering-only (no Wants), so it applies only when gpsd is enabled.
After=network-online.target direwolf.service gpsd.socket gpsd.service
Wants=network-online.target direwolf.service

[Service]
Type=simple
User=renogymon-aprs
Group=renogymon-aprs
EnvironmentFile=-/etc/default/renogymon-aprs
ExecStart=/usr/bin/renogymon-aprs --ssid ${SSID}
Restart=on-failure
RestartSec=10
Environment=RUST_LOG=info
Environment=SSID=N0CALL-1
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
